# daemons.conf
[supervisord]
nodaemon=true

[program:svnserve]
command=/usr/bin/svnserve --root /opt/svn --foreground -d
